# templates/guests/coreos/etcd.service.erb
[Unit]
Description=Clustered etcd
After=etc-discovery.service

[Service]
Restart=always
ExecStart=/usr/bin/etcd -addr localhost:4001 -peer-addr <%= options[:my_ip] %>:7001 -data-dir /home/core/etcd -name etcd-%H -peer-bind-addr=<%= options[:my_ip] %>:7001 -discovery <%= options[:discovery_string] %>

[Install]
WantedBy=local.target
